# deploy/ai-server-scheduler.service
[Unit]
Description=AI Server Scheduler (APScheduler)
After=network.target redis.service
Wants=redis.service

[Service]
Type=exec
WorkingDirectory=/opt/ai-server
EnvironmentFile=/opt/ai-server/.env.prod
ExecStart=/opt/ai-server/.venv/bin/python -m ai_server.scheduler_worker
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-server-scheduler

[Install]
WantedBy=multi-user.target
